# deploy_tools/gunicorn.socket
[Unit]
Description=gunicorn socket

[Socket]
ListenStream=/run/superlists-staging.tobyd.net.sock

[Install]
WantedBy=sockets.target
